# filebeat/package/filebeat_8.6.2-1_armhf/DEBIAN/postinst
#!/bin/sh
systemctl daemon-reload
systemctl enable filebeat
systemctl start filebeat
